# data/configs/blueman-mechanism.service.in
[Unit]
Description=Bluetooth management mechanism

[Service]
Type=dbus
KillMode=process
BusName=org.blueman.Mechanism
ExecStart=@LIBEXECDIR@/blueman-mechanism

[Install]
WantedBy=multi-user.target
