# data/configs/blueman-applet.service.in
[Unit]
Description=Bluetooth management applet

[Service]
Type=dbus
BusName=org.blueman.Applet
ExecStart=@BINDIR@/blueman-applet
