# data/configs/blueman-mechanism.service.in
[Unit]
Description=Bluetooth management mechanism

[Service]
Type=dbus
BusName=org.blueman.Mechanism
ExecStart=@LIBEXECDIR@/blueman-mechanism
