# data/configs/blueman-manager.service.in
[Unit]
Description=Bluetooth Manager

[Service]
Type=dbus
BusName=org.blueman.Manager
ExecStart=@BINDIR@/blueman-manager
